# systemd/dovecot-webadmin-helper.service
[Unit]
Description=Dovecot Web Admin Privileged Helper
After=network.target

[Service]
Type=exec
User=root
# Group=www-data so systemd creates /run/dovecot-webadmin/ as root:www-data
# (with RuntimeDirectoryMode=0750), allowing the web app — which runs as
# www-data with no supplementary groups — to traverse into the dir and
# reach the helper socket. The process still has EUID=0 so all privileged
# subprocesses (useradd, chpasswd, ufw, ...) work normally.
Group=www-data
WorkingDirectory=/opt/dovecot-webadmin
Environment="PATH=/opt/dovecot-webadmin/venv/bin"
EnvironmentFile=-/opt/dovecot-webadmin/.env
ExecStart=/opt/dovecot-webadmin/venv/bin/python -m privileged.server
Restart=always
RestartSec=5

# Runtime directory for socket
RuntimeDirectory=dovecot-webadmin
RuntimeDirectoryMode=0750

# Security hardening — the helper still needs root for useradd/userdel,
# chpasswd, ufw, doveadm, and reading mail logs. Lock down everything else.
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
ProtectProc=invisible
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true

# Filesystem — /usr, /boot, /efi read-only (the helper never writes there
# anyway), but /etc is writable because the helper has to update
# /etc/passwd, /etc/shadow, /etc/ufw/user.rules etc. via useradd / chpasswd
# / ufw. ProtectSystem=strict (which makes /etc read-only too) breaks every
# command that does system administration. ReadWritePaths adds the few
# extra writable dirs ufw and postfix need.
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/var/log /var/spool/postfix /var/mail /home /run/dovecot-webadmin /var/lib/ufw

# Network — Unix socket only, no network sockets needed
RestrictAddressFamilies=AF_UNIX AF_NETLINK

# System call filter — block exotic syscalls that the helper should never use.
# Note: we can't use ~@privileged because in systemd 250+ it includes @chown,
# which the helper legitimately needs at startup to set group ownership on
# /run/dovecot-webadmin/helper.sock so www-data can connect. Instead we list
# the dangerous subgroups individually.
SystemCallFilter=@system-service
SystemCallFilter=~@clock @module @reboot @swap @raw-io @mount @debug @cpu-emulation @obsolete @resources
SystemCallErrorNumber=EPERM
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
