# systemd/dovecot-webadmin.service
[Unit]
Description=Dovecot Web Admin Console
After=network.target dovecot-webadmin-helper.service
Requires=dovecot-webadmin-helper.service

[Service]
Type=exec
User=www-data
Group=www-data
WorkingDirectory=/opt/dovecot-webadmin
Environment="PATH=/opt/dovecot-webadmin/venv/bin"
EnvironmentFile=-/opt/dovecot-webadmin/.env
# --proxy-headers + --forwarded-allow-ips makes uvicorn trust X-Forwarded-For
# ONLY from the local nginx (127.0.0.1), so request.client.host is the real
# client IP. Without it every request looks like 127.0.0.1, which collapses
# the per-IP login rate limit into one global bucket and logs loopback as the
# session IP. Keep the allow-list at 127.0.0.1 — never widen it to '*'.
ExecStart=/opt/dovecot-webadmin/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --proxy-headers --forwarded-allow-ips=127.0.0.1
Restart=always
RestartSec=5

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/dovecot-webadmin/data

# Allow access to helper socket
SupplementaryGroups=www-data

[Install]
WantedBy=multi-user.target
